# docs/docker/cron.d/intranet_crontab
*/5 * * * *       root    /usr/local/bin/php /var/www/html/bin/cron.php
