# docs/docker/cron.d/intranet_crontab
*/5 * * * *       root    php /var/www/html/cron.php
